# domains/templatedata/ipmapper.timer
[Unit]
Description=IPMapper Timer

[Timer]
OnBootSec=10sec
OnUnitActiveSec=3min
Unit=ipmapper.service

[Install]
WantedBy=timers.target
