# domains/templatedata/ipmapper.service
[Unit]
Description=IPMapper Client

[Service]
Type=oneshot
ExecStart=/bin/ipmapper client
